# ops/systemd/trading.service.d/10-runtime-ma-only.conf
[Service]
ExecStart=
ExecStart=/home/ubuntu/trader/trading/.venv/bin/gunicorn -w 1 --timeout 0 --graceful-timeout 900 -b 127.0.0.1:5000 wsgi:application
