# trading/systemd/trading-mem-monitor.service
[Unit]
Description=Trading System Memory and Disk Monitor
After=network-online.target
Wants=network-online.target
# 配置遗漏时限制失败重启频率，避免每 10 秒永久刷日志。
StartLimitIntervalSec=300
StartLimitBurst=3

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/trader/trading
# 最小权限：专用文件只能包含 DINGTALK_WEBHOOK，绝不能换成含 OKX 密钥、
# 登录口令和 FLASK_SECRET_KEY 的 /etc/trading.env。前缀 '-' 允许 config.json
# 已配置 dingtalk.webhook_url 的部署不创建该文件。
EnvironmentFile=-/etc/trading-mem-monitor.env
ExecStart=/home/ubuntu/trader/trading/.venv/bin/python mem_monitor.py
Restart=on-failure
RestartSec=10
UMask=0077
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
